# services/opengarrison-api/deploy/opengarrison-api.service
[Unit]
Description=OpenGarrison API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=opengarrison
Group=opengarrison
WorkingDirectory=/opt/opengarrison-api
Environment=OPENGARRISON_API_DB=/var/lib/opengarrison-api/opengarrison.db
Environment=OPENGARRISON_API_CORS_ORIGINS=https://superganggarrison.com,https://www.superganggarrison.com,https://play.superganggarrison.com,https://unkind-dev.com,https://www.unkind-dev.com
ExecStart=/opt/opengarrison-api/.venv/bin/uvicorn app:app --host 127.0.0.1 --port 8008
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
